# generic-shortcuts_1.0_all/etc/profile.d/generic-shortcuts_1.0_all.sh
. local_shortcuts.dot
. r.path_cleanup.dot
